# systemd-service/FIP_mirror.timer
[Unit]
Description=Run the FIP_mirror script every minute.

[Timer]
OnCalendar=*:*:0/20
Persistent=True

[Install]
WantedBy=timers.target
