# systemd-service/FIP_lastfm_mirror.timer
[Unit]
Description=Run the FIP_lastfm_mirror script every minute.

[Timer]
OnCalendar=*:0/1
Persistent=True

[Install]
WantedBy=timers.target
